# conf/supervisor.conf
[program:__APP__]
directory=__FINALPATH__/app/server/
command=python index.py --port __PORT__
autostart=true
autorestart=true
user=__APP__
