# conf/supervisor.conf
[program:linuxdash]
directory=/var/www/linuxdash/
command=python python-server.py
autostart=true
autorestart=true
