# conf/supervisor.conf
[program:__APP__]
directory=__FINALPATH__/
;command=python app/server/index.py (v2.0)
command=python python-server.py
autostart=true
autorestart=true
user=__USER__
